# crontab
0 0 3 * * * * bash /code/tools/job-runner.sh